# mcp-manager.service
[Unit]
Description=MCP Manager - MCP Server Management Dashboard
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/mcp-manager --port 9847 --config %h/.config/mcp-manager/config.json
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
